# cloud-netconfig/azure/51-cloud-netconfig-hotplug.rules
SUBSYSTEM=="net", KERNEL=="eth*", RUN+="/etc/sysconfig/network/scripts/cloud-netconfig-hotplug"
